# agent/com.fortytwo.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fortytwo.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{AGENT_PATH}}</string>
        <string>--bot-url</string>
        <string>{{BOT_URL}}</string>
        <string>--agent-token</string>
        <string>{{AGENT_TOKEN}}</string>
        <string>--scripts-root</string>
        <string>{{SCRIPTS_ROOT}}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
